# lb-core-env/common/build/hub.dayu-boss.com/outsite/openresty-1.13.6.2/container-files/etc/supervisord.conf
; ----------------------------------------- main -----------------------------------------
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
user=root
logfile=/var/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ----------------------------------------- program -----------------------------------------

;[program:zabbix-agentd]
;command=/usr/sbin/zabbix_agentd -fc /etc/zabbix_agentd.conf
;user=zabbix
;startsecs=5
;startretries=3
;autostart=true
;autorestart=true
;stdout_logfile=/var/logs/supervisor/zabbix-out.log
;stderr_logfile=/var/logs/supervisor/zabbix-err.log


[program:redis-server]
command=/usr/bin/redis-server /etc/redis.conf
user=root
startsecs=1
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

;
[program:crond]
command=/usr/sbin/crond -n
user=root
numprocs=1
startsecs=7
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:cut-ngx-log]
command=/usr/bin/crontab
numprocs=1
user=root
startsecs=8
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

;
[program:lsyncd]
command=/usr/bin/lsyncd -nodaemon /etc/lsyncd.conf
user=root
startsecs=8
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

;
[program:nginx]
command=/usr/local/openresty/bin/openresty -g "daemon off;" -c /usr/local/openresty/nginx/conf/nginx/nginx.conf
user=root
startsecs=10
startretries=3
autostart=true
autorestart=true
stdout_logfile=/dev/null
stderr_logfile=/dev/null

;
[program:monitor_fetch]
command=python /usr/local/openresty/nginx/script/py/MonitorFetch.py
stdout_logfile=/var/logs/monitor_fetch_access.log
stderr_logfile=/var/logs/monitor_fetch_err.log
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true