# lb-core-env/common/build/hub.dayu-boss.com/core/centos-pg10.3/container-files/etc/supervisord.d/postgresqld.ini
[program:postgresqld]
command=/usr/local/postgres/bin/postgres
stdout_logfile=/var/log/pgsql_out.log
stderr_logfile=/var/log/pgsql_err.log
user=postgres
startsecs=8
startretries=3
autostart=true
autorestart=false
