# lb-core-env/common/build/hub.dayu-boss.com/core/openresty-php72/container-files/etc/supervisord.d/nginx.ini
[program:nginx]
command=/usr/local/openresty/bin/openresty -g "daemon off;"
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/crond -n
numprocs=1
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true

[program:cut-ngx-log]
command=/usr/bin/crontab
numprocs=1
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true

[program:php-fpm]
command=/opt/remi/php72/root/usr/sbin/php-fpm --nodaemonize
user=root
startsecs=8
startretries=5
autostart=true
autorestart=true
