# lb-core-env/common/build/hub.dayu-boss.com/core/tengine/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
