# lb-core-env/common/build/hub.dayu-boss.com/core/zookeeper/container-files/etc/supervisord.d/zookeeper.ini
[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground
user=root
startsecs=10
startretries=3
autostart=true
autorestart=true

