# lb-core-env/common/build/hub.dayu-boss.com/core/gb-gather-server/container-files/etc/supervisord.d/gather-server.ini
[program:gather-server]
directory=/usr/local/gather-server/bin/
command=java -classpath ../WEB-INF/classes -Djava.ext.dirs=../WEB-INF/lib org.soul.gather.commons.server.service.GatherServer
user=root
startsecs=6
startretries=3
autostart=true
autorestart=true

