# lb-core-env/common/build/hub.dayu-boss.com/core/tomcat/container-files/etc/supervisord.d/tomcat.ini
[program:tomcat]
directory=/usr/local/tomcat/bin/
command=/usr/local/tomcat/bin/catalina.sh run
user=root
startsecs=2
startretries=3
autostart=true
autorestart=true

[program:crond]
command=/usr/sbin/crond -n
numprocs=1
user=root
startsecs=6
startretries=3
autostart=true
autorestart=true

[program:cleanLOG-15day]
command=/usr/bin/crontab
numprocs=1
stdout_logfile=/dev/null
stderr_logfile=/dev/null
user=root
startsecs=6
startretries=3
autostart=true
autorestart=true
