# lb-core-env/common/build/hub.dayu-boss.com/core/rocketmq/container-files/etc/supervisord.d/rocketmq.ini
[program:mqnamesrv]
command=/usr/local/rocketmq/bin/mqnamesrv
user=root
startsecs=6
startretries=3
autostart=true
autorestart=false

[program:mqbroker]
command=/usr/local/rocketmq/bin/mqbroker -c /usr/local/rocketmq/conf/main/mqbroker.conf
user=root
startsecs=6
startretries=5
autostart=true
autorestart=false

