# lb-core-env/common/build/hub.dayu-boss.com/core/redis/container-files/etc/supervisord.d/redis-server.ini
[program:redis-server]
command=/usr/local/bin/redis-server /etc/redis.conf
user=redis
autostart=true
autorestart=false
startsecs=10
startretries=5

