# lb-core-env/common/build/hub.dayu-boss.com/core/centos-squid/container-files/etc/supervisord.d/squid.ini

[program:squid]
command=/usr/sbin/squid -f /etc/squid/squid.conf -N
user=squid
startsecs=5
startretries=3
autostart=true
autorestart=true
stdout_logfile=/var/log/squid/squid-out.log
stderr_logfile=/var/log/squid/squid-err.log


